# testdata/service-mount/package/missing-fields.service
[Unit]
Description=Service with missing fields