# testdata/service-mount/package/valid.service.in
[Service]
RestartSec=10
User=root
Type=simple
WorkingDirectory=/original-path/
ExecStart=/original-path/bin/executable --argument=/original-path/argument
Restart=always

[Install]
WantedBy=multi-user.target

[Unit]
After=network.target
Description=Valid Service
