# testdata/service-mount/etc/systemd/system/requires.service
[Service]
Type=simple
ExecStart=/package/dir-with-executable/bin/executable --argument=/package/dir-with-executable/argument
Restart=always
RestartSec=10
User=root
WorkingDirectory=/package/dir-with-executable/

[Unit]
After=network.target
Description=Valid Service
Requires=required.service

[Install]
WantedBy=multi-user.target